# ansible/roles/dovecot/templates/var-mail-virtual.automount
[Unit]
Description=Automount Virtual Mail
DefaultDependencies=no

[Automount]
Where={{dovecot_virtual_mailroot}}

[Install]
WantedBy=local-fs.target
