# ansible/roles/certbot/files/certbot-renew.timer
[Unit]
Description=Run Certbot renewal every day
After=network.target

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
