# ansible/roles/mail-nfs-mount/templates/mnt-share-mail.automount
[Unit]
Description=Automount Mail NFS share
DefaultDependencies=no

[Automount]
Where={{mail_nfs_mount_dir}}

[Install]
WantedBy=multi-user.target
